# Gateway/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc.default
user=tor
autorestart=true